# resources/kiwi.upstart
description "IRC web interface"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /var/lib/kiwi/kiwi -f

pre-start exec /var/lib/kiwi/kiwi build
